# dot_config/crontab/kevinnitro
# vim set ft=crontab

# Start tmux as user in order to read environment variables
@reboot /usr/bin/tmux new-session -d
